# systemd/mweigand_eglinfo.service
[Unit]
Description=Mweigand service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/eglinfo

[Install]
WantedBy=multi-user.target
